# com.notion.ical.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.notion.ical.server.plist
  ─────────────────────────────
  Keeps serve_ical.py running in the background so Apple Calendar
  can always reach your .ics files at http://localhost:8080/

  Replace FULL_PATH_TO_SCRIPT and YOUR_MAC_USERNAME before loading.

  Load with:
    cp com.notion.ical.server.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.notion.ical.server.plist
-->
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.notion.ical.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>FULL_PATH_TO_SCRIPT/notion-ical-sync/venv/bin/python3</string>
    <string>FULL_PATH_TO_SCRIPT/notion-ical-sync/serve_ical.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>FULL_PATH_TO_SCRIPT/notion-ical-sync</string>

  <!-- Start immediately when loaded -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Restart automatically if it ever crashes -->
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>FULL_PATH_TO_SCRIPT/Library/Logs/notion-ical-server.log</string>
  <key>StandardErrorPath</key>
  <string>FULL_PATH_TO_SCRIPT/Library/Logs/notion-ical-server-error.log</string>

</dict>
</plist>
